# dist/vm/allspark.service
[Unit]
Description=allspark
After=docker.service

[Service]
StandardOutput=syslog
StandardError=syslog
ExecStart=/allspark/init

[Install]
WantedBy=multi-user.target
