# dist/vm/running-docker/allspark.service
[Unit]
Description=allspark

[Service]
ExecStart=/allspark/init

[Install]
WantedBy=multi-user.target
