# dist/vm/running-docker/allspark.service
[Unit]
Description=allspark
After=docker.service

[Service]
ExecStart=/allspark/init

[Install]
WantedBy=multi-user.target
